# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable appupgrade.service
systemctl start appupgrade.service
