# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop appupgrade.service
systemctl disable appupgrade.service
